# layout/rootfs/run/service/s6-linux-init-early-getty/run
#!/command/execlineb -P

getty -L -H %%HOSTNAME%% 38400 /dev/%%CONSOLE%%
