# layout/rootfs/etc/s6-linux-init/run-image/service/s6-svscan-log/run
#!/command/execlineb -P

redirfd -w 2 /dev/console
redirfd -w 1 /dev/null
redirfd -rnb 0 fifo
s6-applyuidgid -u 98 -g 98 --
s6-log -bp -- t /run/uncaught-logs
